# initscripts/init.d/squashfsmount
#!/sbin/runscript

description="Mounts squashfs image files"

depend(){
	need localmount
}

start(){
	ebegin "Mounting squashfs image filesystems"
	mount -t squashfs -O ro,defaults /var/squashfs/portage/gentoo.sqfs "${portage_dir}"
	mount -t squashfs -O ro,defaults /var/squashfs/kernel/kernel.sqfs "${kernel_dir}"
	eend $?
}

stop(){
	ebegin "Unmounting squashfs image filesystems"
	umount "${portage_dir}"
	umount "${kernel_dir}"
	eend $?
}
